# config/systemd/user/redshift.service
[Unit]
Description=Redshift Daemon
After=xorg.target

[Service]
ExecStart=/usr/bin/redshift -l 45.61:8.9

[Install]
WantedBy=mystuff.target
