# config/systemd/user/xwallpaper.timer
[Timer]
OnUnitInactiveSec=900
Unit=xwallpaper.service

[Install]
WantedBy=mystuff.target
