# xorg/.config/systemd/user/xlogin.target
[Unit]
Description=X login target
Wants=display-manager.service
AllowIsolate=yes
