# xorg/.config/systemd/user/xlogin.target
[Unit]
Description=X login target
AllowIsolate=yes
