# config/systemd/custom/user@.service
.include /usr/lib/systemd/system/user@.service
[Service]
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/$I
TimeoutStopSec=2
