# config/autostart/compton.desktop
[Desktop Entry]
Type=Application
TryExec=compton
Exec=compton --config ~/.config/compton.conf
Name=Compositor
