# config/autostart/keyboard.desktop
[Desktop Entry]
Type=Application
TryExec=setxkbmap
Exec=setxkbmap it
Name=Keyboard Setter
