# config/autostart/nitrogen.desktop
[Desktop Entry]
Type=Application
TryExec=nitrogen
Exec=nitrogen --restore
Name=Wallpaper setter
