# config/autostart/numlock.desktop
[Desktop Entry]
Type=Application
TryExec=numlockx
Exec=numlockx on
Name=Keyboard Setter
