# config/autostart/thunar.desktop
[Desktop Entry]
Type=Application
TryExec=thunar
Exec=thunar --daemon
Name=File Manager
