# config/autostart/urxvtd.desktop
[Desktop Entry]
Type=Application
TryExec=numlockx
Exec=urxvtd -q -o
Name=Urxvtd
