# config/autostart/urxvtd.desktop
[Desktop Entry]
Type=Application
TryExec=urxvtd
Exec=urxvtd -q -o
Name=Urxvtd
