# config/autostart/xrdb.desktop
[Desktop Entry]
Type=Application
TryExec=xrdb
Exec=xrdb -merge ~/.config/Xresources
Name=Xrdb
